# supervisord.conf
[program:node]
command=node vanity.js
environment=NODE_ENV=dotcloud
directory=/home/dotcloud/current
